# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:flask_app]
command=gunicorn -c gunicorn.conf.py app:app
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/flask_app.err.log
stdout_logfile=/var/log/flask_app.out.log

[program:celery_worker]
command=python start_worker.py
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_worker.err.log
stdout_logfile=/var/log/celery_worker.out.log

[program:celery_beat]
command=python start_beat.py
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/celery_beat.err.log
stdout_logfile=/var/log/celery_beat.out.log

[program:redis]
command=redis-server
autostart=true
autorestart=true
stderr_logfile=/var/log/redis.err.log
stdout_logfile=/var/log/redis.out.log